# openstick-gc-manager/adbd.service
[Unit]
Description=ADB Daemon

[Service]
Type=simple
ExecStart=/usr/bin/adbd
Restart=on-failure

[Install]

